# deploy/harbor.service
# systemd unit for running harbor as a daemon.
#
# Place this at /etc/systemd/system/harbor.service and:
#
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now harbor
#   sudo journalctl -u harbor -f       # follow the log
#
# Customize the [Service] env vars below for your deployment. The
# bootstrap SQL file (deploy/harbor-bootstrap.sql) is what actually
# loads the extension and starts the server — edit it to match your
# auth + CORS + bind config before enabling the unit.

[Unit]
Description=harbor — DuckDB as an HTTP service
Documentation=https://github.com/shreeve/duckdb-harbor
After=network-online.target
Wants=network-online.target

[Service]
# A non-root user; create with:  sudo useradd --system --create-home --home-dir /var/lib/harbor harbor
User=harbor
Group=harbor

# Where the bootstrap SQL + extension binary + database file live.
WorkingDirectory=/var/lib/harbor
Environment=DUCKDB_BIN=/usr/local/bin/duckdb
Environment=HARBOR_BOOTSTRAP=/var/lib/harbor/harbor-bootstrap.sql

# Foreground duckdb process. The bootstrap SQL ends with `CALL harbor_wait();`
# which blocks until SIGTERM, so this stays alive under systemd.
#
# `-no-stdin` is critical: without it, duckdb's REPL would race with
# systemd-journald for stdin and exit immediately.
ExecStart=/bin/sh -c '${DUCKDB_BIN} -unsigned -no-stdin < ${HARBOR_BOOTSTRAP}'

# Graceful shutdown:
#   - harbor's signal handler treats SIGTERM as "stop the server cleanly"
#   - 30s gives in-flight queries a chance to finish; longer queries
#     get interrupted via the same path that harbor_query_timeout_s uses
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=5

# Hardening — none of these are strictly required, all of them are good
# practice. Comment out individual lines if your deployment needs e.g.
# write access to a wider filesystem.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/harbor
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=@system-service

# A loopback bind is the default; if you bind harbor on an external
# interface, put a real reverse proxy (nginx, caddy, traefik) in front
# for TLS. See docs/DEPLOYMENT.md.

[Install]
WantedBy=multi-user.target
